# cron-backup/crontab
0 1 * * * root /usr/bin/backup.sh >> /var/log/cron.log 2>&1
